# s6/services.d/proxysql/run
#!/usr/bin/execlineb -P


proxysql -f -M -D /var/lib/proxysql
